# provision/files/rc.local
#!/bin/bash 

cd /opt && docker-compose up -d
exit 0

